# scripts/crontab
# START CRON JOB
*/1 * * * * /usr/app/scripts/run-main.sh
# END CRON JOB
